# .config/systemd/user/lippy-sync-data.timer
[Unit]
Description=Backup items to aws bucket altinners-creds

[Timer]
OnBootSec=1min
OnUnitActiveSec=1day

[Install]
WantedBy=timers.target
