# .config/systemd/user/lippy-sync-data.service
[Service]
ExecStart=/home/andrew.lipscomb/bin/lippy-sync-data
Type=oneshot
