# .config/systemd/user/lippys-pc.service
[Service]
ExecStart=/bin/bash /home/andrew.lipscomb/bin/get-comp-back

[Install]
WantedBy=default.target
