# .config/systemd/user/picom.service
[Unit]
Description=Picom

[Service]
ExecStart=/bin/bash -c /usr/bin/picom
