# .config/systemd/user/swarmfarm-auth-refresh.timer
[Unit]
Description=Refresh SwarmFarm user tokens periodically

[Timer]
OnBootSec=1min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target

