# .config/systemd/user/swarmfarm-wait-for-robot@.service
[Unit]
Description=Waits for robot to come on

[Service]
ExecStart=%h/bin/swarmfarm-wait-for-robot %i
