# .config/systemd/user/syncthing.service
[Unit]
Description=Syncthing agent

[Service]
Type=simple
ExecStart=/usr/bin/syncthing serve --no-browser

[Install]
WantedBy=default.target
